# scripts/flight_prediction/flight_prediction_agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.you.flight_prediciton_agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3'</string>
        <string>/Users/user/automate/flight_prediction.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/user/Library/Logs/flight_prediciton_agent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/user/Library/Logs/flight_prediciton_agent-error.log</string>

</dict>
</plist>